# dist/linux/tf2richpresence.service
[Unit]
Description=TF2 Rich Presence
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /usr/share/tf2-rich-presence/src/main.py

[Install]
WantedBy=default.target