# dist/linux/tf2richpresence.service
[Unit]
Description=TF2 Rich Presence
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/env python3 /usr/share/tf2-rich-presence/src/main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target