# dist/linux/tf2richpresence.service
[Unit]
Description=TF2 Rich Presence
Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/usr/share/tf2-rich-presence/open_tf2_rich_presence.sh
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target